# debian/cron.d
MAILTO=chilts@appsattic.com

# mins hour dofm mnth dofw   user          command
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
   */5    *    *    *    *   root          dropmail
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
